# 00-Setup/Ansible/roles/filesystem/files/vol.mount
[Unit]
Description=Mount vol to /vol
DefaultDependencies=no

[Mount]
What=/dev/vdb1
Where=/vol
Type=ext4
